# gitbook-task.cron
crontab<<EOF
*/3 * * * * sh /home/dlm/work/mango-blog/gitbook-execute.sh
EOF